# SETUP/files/etc/init.d/fhem
#!/bin/sh /etc/rc.common
# FHEM Init Script
START=99
STOP=15
start() {
    pl /opt/fhem/fhem.pl /opt/fhem/fhem.cfg
}
stop(){
    echo &quot;shutdown&quot; | nc localhost 7072
}
